# packages/coding-agent/systemd/pi-architect.service
[Unit]
Description=Pi resident architect
After=default.target

[Service]
Type=exec
WorkingDirectory=%h
ExecStart=@PI_NODE_LAUNCHER@ --tsconfig @PI_TSCONFIG@ @PI_CLI_SOURCE@ architect
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
